# docker/rootfs/etc/s6-overlay/s6-rc.d/snmp-bridge/run
#!/command/with-contenv sh
exec /app/snmp-bridge
